# deployment/go2rtc.service
[Unit]
Description=go2rtc RTSP to HLS/WebRTC Streaming Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/cctv/go2rtc
ExecStart=/var/www/cctv/go2rtc/go2rtc -config /var/www/cctv/go2rtc/go2rtc.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
